# cookbooks/mysql/files/default/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

echo noop > /sys/block/sdb/queue/scheduler
echo 256 > /sys/block/sdb/queue/nr_requests
echo 1 > /sys/block/sdb/device/queue_depth

exit 0
